# docker/macdiscovery.service
[Unit]
Description=Service to manage firmware container for robot
After=docker.service
Requires=docker.service

[Service]
Type=forking
#Restart=always
WorkingDirectory=/home/pi/git/gritsbot_2/docker/
ExecStartPre=-/usr/bin/docker stop firmware
ExecStartPre=-/usr/bin/docker rm firmware
ExecStartPre=-/usr/bin/docker pull robotarium/firmware
#ExecStartPre=/bin/bash -c "/bin/systemctl set-environment SERIAL_PORT=$(sudo -u pi python3 detect_serial.py)"
#ExecStart=/usr/bin/docker run -d --restart always --name firmware --net host --device ${SERIAL_PORT}:/dev/ttyACM0 robotarium/firmware
ExecStart=/bin/bash docker_run.sh
#ExecStop=-/usr/bin/docker stop firmware
#ExecStop=/usr/bin/docker rm -f firmware
ExecReload=/usr/bin/docker restart firmware

[Install]
WantedBy=multi-user.target
